# local.claude-schedule.service.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Template for the claude-schedule-management web service launch agent.
  Replace __ROOT__ and __NODE__ before installing.
  See README.md for the install command.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.claude-schedule.service</string>
  <key>ProgramArguments</key>
  <array>
    <string>__NODE__</string>
    <string>__ROOT__/server/dist/index.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__ROOT__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PORT</key>
    <string>7878</string>
    <key>HOST</key>
    <string>127.0.0.1</string>
    <key>PATH</key>
    <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__ROOT__/logs/service-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__ROOT__/logs/service-stderr.log</string>
</dict>
</plist>
